# deploy/systemd/workai-assess.timer
[Unit]
Description=Schedule WorkAI v2 assess daily run

[Timer]
OnCalendar=*:6/15
Persistent=true
Unit=workai-assess.service

[Install]
WantedBy=timers.target
