# deploy/systemd/workai-normalize.service
[Unit]
Description=WorkAI v2 normalize run
After=network.target

[Service]
Type=oneshot
User=workai
Group=workai
WorkingDirectory=/opt/workai
EnvironmentFile=/etc/workai/secrets/db.env
EnvironmentFile=/etc/workai/secrets/workai.env
ExecStart=/opt/workai/.venv/bin/python /opt/workai/scripts/workai_normalize.py run
