# deploy/systemd/workai-verify-units.timer
[Unit]
Description=Schedule WorkAI v2 verify units

[Timer]
OnCalendar=hourly
Persistent=true
Unit=workai-verify-units.service

[Install]
WantedBy=timers.target
