# deploy/systemd/workai-audit-manual.service
[Unit]
Description=WorkAI v2 manual AI audit run
After=network.target

[Service]
Type=oneshot
User=workai
Group=workai
WorkingDirectory=/opt/workai
EnvironmentFile=/etc/workai/secrets/db.env
EnvironmentFile=/etc/workai/secrets/api.env
EnvironmentFile=/etc/workai/secrets/workai.env
Environment=HOME=/opt/workai
PassEnvironment=WORKAI_AUDIT_EMPLOYEE_ID WORKAI_AUDIT_TARGET_DATE WORKAI_AUDIT_FORCE
ExecStart=/opt/workai/scripts/run_audit_prod.sh
