# deploy/systemd/workai-assess.timer
[Unit]
Description=Schedule WorkAI v2 assess daily run

[Timer]
OnCalendar=*-*-* 01:15:00
Persistent=true
Unit=workai-assess.service

[Install]
WantedBy=timers.target
